# shell/test/crontab
*/1 * * * * /home/yoshinoakira/prog/shell/test/autoJava
